# nenuscanner.service
[Unit]
Description=NenuScanner
After=network.target

[Service]
RestartSec=2s
Type=simple
WorkingDirectory=/home/pi/nenuscanner
ExecStart=python -m server
Restart=always

[Install]
WantedBy=default.target

