# nenuscanner.service
[Unit]
Description=NenuScanner
After=network.target

[Service]
RestartSec=2s
Type=simple
WorkingDirectory=/home/pi/nenuscanner/server
ExecStart=flask --app . run --host 0.0.0.0
Restart=always

[Install]
WantedBy=default.target

