# nenuscanner.service
[Unit]
Description=NenuScanner
Wants=mediamtx.service
After=network.target mediamtx.service

[Service]
RestartSec=2s
Type=simple
WorkingDirectory=/home/pi/nenuscanner/server
ExecStart=/home/pi/.local/bin/poetry run nenuscanner
Restart=always

[Install]
WantedBy=default.target

